# tag-laptop/config/systemd/user/compton.service
[Unit]
Description=compton

[Service]
ExecStart=/usr/bin/compton --config %h/.config/compton/compton.conf
