# tag-laptop/config/systemd/user/polybar.service
[Unit]
Description=Start polybar

[Service]
Restart=on-abnormal
ExecStart=/usr/sbin/polybar laptop
ExecReload=/usr/sbin/polybar -r laptop
