# tag-laptop/config/systemd/user/wallpaper.service
[Unit]
Description=Setup wallpaper

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/feh --bg-fill %h/.cache/wallpaper/wallpaper
